# packages/composer-ui/supervisord-concerto-connector-server.conf
[program:concerto-connector-server]
command=concerto-connector-server
user=concerto
environment=HOME="/home/concerto",USER="concerto"
