# packages/composer-ui/supervisord-composer-connector-server.conf
[program:composer-connector-server]
command=concerto-connector-server
user=composer
environment=HOME="/home/composer",USER="composer"
