# packages/composer-playground/supervisord-nginx.conf
[program:nginx]
command=nginx -g "daemon off;"
